# docker/supervisord/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/www/html/src/storage/logs/supervisord.log
logfile_maxbytes=50MB
logfile_backups=3
loglevel=info
pidfile=/var/www/html/src/storage/logs/supervisord.pid
childlogdir=/var/www/html/src/storage/logs

[program:php-fpm]
command = /usr/local/sbin/php-fpm --nodaemonize
autostart = true
autorestart = true
stderr_logfile = /var/www/html/src/storage/logs/php-fpm.err.log
stdout_logfile = /var/www/html/src/storage/logs/php-fpm.out.log
stderr_logfile_maxbytes = 10MB
stdout_logfile_maxbytes = 10MB
stderr_logfile_backups = 5
stdout_logfile_backups = 5
